# config/launchagents/ai.ganchitecture.lisbon-esp32-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- ESP32 soundscape sync. Reads the SWN bridge status JSON and drives
     the dual-strip dystopia firmware via /dev/cu.usbserial-0001.
     Depends on the SWN bridge writing status JSON. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.ganchitecture.lisbon-esp32-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ganchitecture/code/lisbon-av-install/lighting/run_esp32_soundscape_sync.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ganchitecture/code/lisbon-av-install</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/ganchitecture/Library/Logs/lisbon-esp32-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ganchitecture/Library/Logs/lisbon-esp32-sync.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
