# scripts/launchagents/ai.ganchitecture.lisbon-camera-probe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.ganchitecture.lisbon-camera-probe</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/open</string>
    <string>-W</string>
    <string>/Users/ganchitecture/code/lisbon-av-install/cv/camera_probe/LisbonCameraProbe.app</string>
    <string>--args</string>
    <string>--snapshot-path</string>
    <string>/Users/ganchitecture/code/lisbon-av-install/cv/captures/latest.jpg</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/ganchitecture/Library/Logs/lisbon-camera-probe.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ganchitecture/Library/Logs/lisbon-camera-probe.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
